# .config/syncuser/dotfiles.list
# Dotfiles: Copy only if missing

$HOME/.bash_profile
$HOME/.bashrc
$HOME/.bashrc.d/
$HOME/.gitconfig
$HOME/.ssh/config
$HOME/.editorconfig
$HOME/.harlequin.toml
